# archives/system-snapshots/openclaw-archive/scripts/reed-dull-proactive-pulse.timer
[Unit]
Description=REED DULL — Proactive Pulse Timer (every 6 hours)

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
AccuracySec=1min
Persistent=true
Unit=reed-dull-proactive-pulse.service

[Install]
WantedBy=timers.target
