# archives/system-snapshots/openclaw-archive/scripts/reed-dull-weekly-gap.timer
[Unit]
Description=REED DULL — Weekly Knowledge Gap Timer (Friday 17:00 WIB)

[Timer]
OnCalendar=Fri *-*-* 17:00:00
AccuracySec=1min
Persistent=true
Unit=reed-dull-weekly-gap.service

[Install]
WantedBy=timers.target
