# archives/system-snapshots/openclaw-archive/scripts/reed-dull-weekly-gap.service
[Unit]
Description=REED DULL — Weekly Knowledge Gap Report (Friday 17:00 WIB)

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/openclaw/banirisset/scheduler/scripts/proactive_engine.py --mode weekly_gap
WorkingDirectory=/home/openclaw/banirisset
EnvironmentFile=/home/openclaw/.openclaw/openclaw-gateway.env

[Install]
WantedBy=default.target
